# com.user.gdrive-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.gdrive-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Pfad/zu/sync_gdrive_to_proton.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer> <!-- 1 Stunde -->
</dict>
</plist>
